# runit/mailer.sender/run
#!/bin/sh

exec 2>&1

cd /www/mailer

exec php send-daemon.php
